# nzbget.conf
[program:nzbget-precheck]
user = root
command = /tmp/nzbget-precheck.sh

[program:nzbget]
user = nobody
command = /usr/bin/nzbget -D -c /config/nzbget.conf