# nzbget.conf
[program:nzbget]
user = nobody
command = /usr/bin/nzbget -D -c /config/nzbget.conf