# nzbget.conf
[program:start]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/start.sh

[program:nzbget]
user = nobody
command = /usr/bin/nzbget -D -c /config/nzbget.conf
umask = 000