# nzbget.conf
[program:prerun]
user = root
command = /etc/supervisor/conf.d/prerun.sh

[program:nzbget]
user = nobody
command = /usr/bin/nzbget -D -c /config/nzbget.conf