# setup/nzbget.conf
[program:nzbget]
user = nobody
command = /home/nobody/setup.sh
umask = 000